# cm4/config/cerebrum-tunnel.service
[Unit]
Description=Persistent reverse SSH tunnel for Cerebrum (CM4 -> VPS)
After=network-online.target
Wants=network-online.target

[Service]
User=kali
Environment="AUTOSSH_GATETIME=0"
ExecStart=/usr/bin/autossh -M 0 -N \
  -o "ServerAliveInterval 30" -o "ServerAliveCountMax 3" \
  -o "ExitOnForwardFailure=yes" \
  -i /home/kali/.ssh/cerebrum_tunnel_id_ed25519 \
  -R 127.0.0.1:9000:127.0.0.1:7000 unicorn1@YOUR_VPS_IP_HERE
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
